# bot/google-cloud-bpsr-bot.service.example
[Unit]
Description=BPSR Discord Recruitment Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/bpsr/bpsr
ExecStart=/usr/bin/npm run bot:start
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target